# sn9trader.cron
0 2 * * * /Users/samarth/Downloads/Official/PROJECTS/CODES/Python/Projects/sn9trader/venv/bin/python3 /Users/samarth/Downloads/Official/PROJECTS/CODES/Python/Projects/sn9trader/congress_tracker.py >> /Users/samarth/Downloads/Official/PROJECTS/CODES/Python/Projects/sn9trader/logs/congress.log 2>&1
0 */6 * * * /Users/samarth/Downloads/Official/PROJECTS/CODES/Python/Projects/sn9trader/venv/bin/python3 /Users/samarth/Downloads/Official/PROJECTS/CODES/Python/Projects/sn9trader/ceo_tracker.py >> /Users/samarth/Downloads/Official/PROJECTS/CODES/Python/Projects/sn9trader/logs/ceo.log 2>&1
*/5 * * * * /Users/samarth/Downloads/Official/PROJECTS/CODES/Python/Projects/sn9trader/venv/bin/python3 /Users/samarth/Downloads/Official/PROJECTS/CODES/Python/Projects/sn9trader/whale_tracker.py >> /Users/samarth/Downloads/Official/PROJECTS/CODES/Python/Projects/sn9trader/logs/whale.log 2>&1
*/15 * * * * /Users/samarth/Downloads/Official/PROJECTS/CODES/Python/Projects/sn9trader/venv/bin/python3 /Users/samarth/Downloads/Official/PROJECTS/CODES/Python/Projects/sn9trader/options_tracker.py >> /Users/samarth/Downloads/Official/PROJECTS/CODES/Python/Projects/sn9trader/logs/options.log 2>&1
0 */6 * * * /Users/samarth/Downloads/Official/PROJECTS/CODES/Python/Projects/sn9trader/venv/bin/python3 /Users/samarth/Downloads/Official/PROJECTS/CODES/Python/Projects/sn9trader/signal_combiner.py >> /Users/samarth/Downloads/Official/PROJECTS/CODES/Python/Projects/sn9trader/logs/signals.log 2>&1
0 22 * * * /Users/samarth/Downloads/Official/PROJECTS/CODES/Python/Projects/sn9trader/venv/bin/python3 /Users/samarth/Downloads/Official/PROJECTS/CODES/Python/Projects/sn9trader/paper_trader.py >> /Users/samarth/Downloads/Official/PROJECTS/CODES/Python/Projects/sn9trader/logs/paper.log 2>&1
