# logs/systemd_temp/sn9trader-whale.timer
[Unit]
Description=Sn9Trader Whale Transaction Tracker Timer

[Timer]
OnCalendar=*:0/5
Persistent=true

[Install]
WantedBy=timers.target
