# logs/systemd_temp/sn9trader-congress.timer
[Unit]
Description=Sn9Trader Congressional Trades Tracker Timer

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
