# logs/systemd_temp/sn9trader-paper.timer
[Unit]
Description=Sn9Trader Paper Portfolio Tracker Timer

[Timer]
OnCalendar=*-*-* 22:00:00
Persistent=true

[Install]
WantedBy=timers.target
